# admin-bot.service
[Unit]
Description=Admin Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpn-service/services
Environment=PATH=/opt/vpn-service/services/venv/bin
ExecStart=/opt/vpn-service/services/venv/bin/python start_admin_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 