# services/vpnbot.target
[Unit]
Description=VPN Bot Complete System
Documentation=https://github.com/your-repo/vpnbot
After=network.target
Wants=network.target

[Install]
WantedBy=multi-user.target 